# systemd/user/wlsctx@.socket
[Unit]
Description=Wayland socket with security context attached (%i)
PartOf=graphical-session.target
After=graphical-session.target
Conflicts=wayland-session-shutdown.target
Before=wayland-session-shutdown.target

[Socket]
RuntimeDirectory=wlsctx/%i
RuntimeDirectoryPreserve=yes
RuntimeDirectoryMode=0700
SocketMode=0600
ListenStream=%t/wlsctx/%i/wayland-1

[Install]
WantedBy=%i.target
