# containers/systemd/wl-xwayland@.container
[Unit]
Description=Sandboxed Xwayland satellite
After=wl-display@.container
Wants=wl-display@.container
StopWhenUnneeded=yes

[Service]
ExecStartPost=-podman exec %p_%i xhost SI:localuser:%u

[Container]
Exec=xwayland-satellite
Image=xwayland-satellite.image
Pull=never
RunInit=true
AddDevice=-/dev/dri/renderD128
AddDevice=-/dev/dri/renderD129
AddDevice=-/dev/dri/renderD130
AddDevice=-/dev/dri/renderD131
PodmanArgs=--ipc=host
Volume=wl-%j_cache:%h/.cache
PodmanArgs=--hostuser=%u
